# service/redis/run
#!/bin/sh
echo "maxmemory 50mb" > /etc/redis.conf
exec /usr/bin/redis-server /etc/redis.conf
