# service/redis/run
#!/bin/sh
sed -i -e '/^# maxmemory/a\maxmemory 50mb' /etc/redis/redis.conf
exec /usr/bin/redis-server
