# service/rspamd/run
#!/bin/sh
exec /usr/bin/rspamd -u _rspamd -g _rspamd -f

