# service/unbound/run
#!/bin/sh
exec /usr/sbin/unbound -d -c /etc/unbound.conf
